# scripts/fakesmtp.service
[Unit]
Description=MailCatch Email Testing Server
Documentation=https://github.com/your-username/mailcatch
After=network.target
Wants=network.target

[Service]
Type=simple
User=mailcatch
Group=mailcatch
Environment=SMTP_PORT=2525
Environment=HTTP_PORT=8080
Environment=CLEAR_ON_SHUTDOWN=true
ExecStart=/usr/local/bin/mailcatch --daemon
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
TimeoutStopSec=30

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/mailcatch

# Resource limits
LimitNOFILE=65536
LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target